# playbooks/roles/provision/templates/etc/update-motd.d/01-custom.j2
#!/bin/sh

{{ ansible_managed | comment }}

{{ motd_body }}