# ansible/playbooks/roles/tanium_client/templates/etc/systemd/system/taniumclient.service.d/90-cgroup.conf.j2

{{ ansible_managed | comment }}

[Service]
Slice=taniumclient.slice
CPUAccounting=true
MemoryAccounting=true
TasksAccounting=true
IOAccounting=true
IPAccounting=true
